# deploy/pi/walle.service
# ~/.config/systemd/user/walle.service
# WALL-E AI — user-level systemd service for Raspberry Pi 4 / 5
# Created by K.Astra and its members.
#
# Installed automatically by deploy/pi/install.sh
# To edit manually: nano ~/.config/systemd/user/walle.service
#                   systemctl --user daemon-reload && systemctl --user restart walle

[Unit]
Description=WALL-E AI Voice Assistant Server
Documentation=https://github.com/kastra/walle-ai
After=network-online.target sound.target
Wants=network-online.target
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__

# Load environment from .env file
EnvironmentFile=__PROJECT_DIR__/.env

# Set PYTHONUNBUFFERED so logs appear in journalctl in real time
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# Use the virtual environment's Python
ExecStart=__VENV_DIR__/bin/python main.py

# Restart automatically if it crashes (e.g., network blip or Gemini quota)
Restart=on-failure
RestartSec=5s
# Give it 30s to gracefully shut down before SIGKILL
TimeoutStopSec=30

# Resource limits — optimized for Pi 4 (4GB RAM / CPU prioritization)
# Cap memory at 3GB to leave room for OS and audio servers (PulseAudio/PipeWire)
MemoryMax=3G
# CPU weight (100 = normal priority)
CPUWeight=100

# Security hardening (user-level services support basic sandbox)
NoNewPrivileges=true
PrivateTmp=true

# Log to walle.log file inside the project directory
StandardOutput=append:__PROJECT_DIR__/walle.log
StandardError=append:__PROJECT_DIR__/walle.log
SyslogIdentifier=walle-ai

[Install]
WantedBy=default.target
